# telegram-bot-cleanup.service.template
[Unit]
Description=Claude Telegram Bot - очистка временных файлов

[Service]
Type=oneshot
User={{USER}}
WorkingDirectory={{BOT_DIR}}
EnvironmentFile={{BOT_DIR}}/.env
ExecStart={{NODE_BIN}} {{BOT_DIR}}/cleanup.js
