# telegram-bot-cleanup.timer.template
[Unit]
Description=Ежедневная очистка временных файлов Claude Telegram Bot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
